# files/etc/systemd/system/neutron-openvswitch-agent.service.d/customize.conf
[Service]
Restart=always
KillMode=control-group
